# src/modules/adsb-feeder/filesystem/root/usr/lib/systemd/system/adsb-nonimage.service
[Unit]
Description=ADS-B Feeder Dummy service

[Service]
Type=oneshot
WorkingDirectory=/opt/adsb
ExecStart=/bin/bash -c "exit 0"
TimeoutStartSec=0

