# src/modules/adsb-feeder/filesystem/root/usr/lib/systemd/system/adsb-docker.service
[Unit]
Description=ADSB Docker Service
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/adsb
ExecStartPre=/usr/bin/bash /opt/adsb/.setup-env.sh
ExecStart=/usr/bin/docker-compose-start
ExecStop=/usr/bin/docker compose down -t 30
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
