# src/modules/adsb-feeder/filesystem/root/usr/lib/systemd/system/adsb-netdog.service
[Unit]
Description=ADS-B Feeder Image Network Watchdog
After=network.target

[Service]
WorkingDirectory=/opt/adsb
ExecStart=/opt/adsb/scripts/netdog.sh
SyslogIdentifier=adsb-netdog
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
