# src/modules/adsb-pi-setup/filesystem/root/usr/lib/systemd/system/adsb-setup.service
[Unit]
Description=ADSB Setup Service
After=network-online.target

[Service]
RemainAfterExit=yes
WorkingDirectory=/usr/local/share/adsb-pi-setup
# ppa deadsnakes python 3.11
ExecStart=/usr/bin/python3 /usr/local/share/adsb-pi-setup/app.py
Restart=always

[Install]
WantedBy=multi-user.target
