# src/tools/automated-boot-testing/adsb-test-service.service
[Unit]
Description=ADS-B Feeder Test Service
Documentation=https://github.com/dirkhh/adsb-feeder-image
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/adsb-test-service
ExecStart=/opt/adsb-test-service/venv/bin/python /opt/adsb-test-service/adsb-test-service.py --config /etc/adsb-test-service/config.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=adsb-test-service

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/etc/adsb-test-service /tmp /srv/iscsi /srv/tftp /etc/tgt /mnt /opt/adsb-test-service/test-images /opt/adsb-test-service/serial-logs /home/testuser /var/lib/adsb-test-service

# Environment
Environment=PYTHONPATH=/opt/adsb-test-service
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
