# src/modules/adsb-feeder/filesystem/root/usr/lib/systemd/system/adsb-docker.service.d/adsb-docker.conf
# this sets the variables used in the service file
# we need two because we cannot embed a space (seriously)
# set this to FR24_FLAG=-f  and FR24_VALUE=fr24.yml
[Service]
Environment=FR24_FLAG=
Environment=FR24_VALUE=
